# upstart/muzzcar.conf
# muzzcar upstart script
# To be saved as /etc/init/muzzcar.conf
#

description     "Muzzley MuzzCar Upstart job"
author          "Tiago Alves <tiago.alves@muzzley.com>"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [06]

# Prepare the environment
#   Create directories for logging and process management
#   Change ownership to the user running the process
pre-start script
    mkdir -p /var/opt/muzzcar
    chown -R muzzley:muzzley /var/opt/muzzcar
end script

# If the process quits unexpectedly trigger a respawn
respawn

# Start the process
exec start-stop-daemon --start --chuid muzzley --chdir /home/pi/muzzley/muzzcar --make-pidfile --pidfile /var/opt/muzzcar/run/muzzcar.pid --exec /usr/bin/node -- /home/pi/muzzley/muzzcar >> /var/opt/muzzcar/log/muzzcar.log 2>&1
